# docker/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /Savior/docker/savior.ini

[program:nginx-app]
command = /usr/sbin/nginx